# legacy/fleet/systemd/openclaw-fleet-dashboard.service
[Unit]
Description=OpenClaw fleet dashboard aggregation
Wants=network-online.target
After=network-online.target default.target

[Service]
Type=oneshot
EnvironmentFile=__CONFIG_FILE__
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/bin/openclaw-fleet-dashboard-run.sh --config __CONFIG_FILE__
