# systemd/openclawnurse.service
[Unit]
Description=OpenClawNurse daily doctor run
After=default.target

[Service]
Type=oneshot
EnvironmentFile=__CONFIG_FILE__
EnvironmentFile=-%h/.openclaw/gateway.systemd.env
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/bin/openclaw-doctor.sh --config __CONFIG_FILE__
MemoryMax=1200M
MemorySwapMax=256M
CPUQuota=60%
TasksMax=128
