# systemd/openclawnurse.timer
[Unit]
Description=Run OpenClawNurse on schedule

[Timer]
OnCalendar=__ON_CALENDAR__
Persistent=true
Unit=openclawnurse.service

[Install]
WantedBy=timers.target
